# systemd/rotrocks-scraper.timer
[Unit]
Description=Run RotDotRocks price scraper every 6h

[Timer]
# Fires at 00:00, 06:00, 12:00, 18:00 UTC every day
OnCalendar=*-*-* 00,06,12,18:00:00
# Spread invocations to avoid hammering Eldorado at exactly :00:00
RandomizedDelaySec=300
# If the box was off when a fire-time passed, run on next boot
Persistent=true
Unit=rotrocks-scraper.service

[Install]
WantedBy=timers.target
